# wthr_srvc.timer
[Unit]
Description=timer that updates the weather hourly 

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
